# suitecrm-scheduler@.timer
[Unit]
Description=SuiteCRM %i scheduler

[Timer]
OnUnitInactiveSec=1min

[Install]
WantedBy=multi-user.target
