# suitecrm-scheduler@.timer
[Unit]
Description=SuiteCRM %i scheduler timer

[Timer]
OnBootSec=1min
OnUnitInactiveSec=1min

[Install]
WantedBy=multi-user.target
