# suitecrm-scheduler.timer
[Unit]
Description=Run SuiteCRM scheduler

[Timer]
OnUnitInactiveSec=1min
Unit=suitecrm-scheduler.service

[Install]
WantedBy=multi-user.target
